# container_scripts/services/starbound/run
#!/bin/bash
set -e

/data/steam/starbound_server/linux64/starbound_server
